# recipes-nms/nms-server/files/crontab
@reboot systemctl restart systemd-networkd
